# tools/systemd/folkrnn-backup.timer
[Unit]
Description = Run backup action folk_rnn.org daily

[Timer]
OnBootSec=15min
OnCalendar=05:43

[Install]
WantedBy = timers.target