# tools/systemd/letsencrypt.timer
[Unit]
Description = Run Lets Encrypt check daily

[Timer]
OnCalendar=04:32

[Install]
WantedBy = timers.target